# deploy_tools/templates/gunicorn.systemd.template.socket.j2
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/{{ domain }}.sock

[Install]
WantedBy=sockets.target
